# templates/default/repo-sync.cron
0 0 * * * /root/repo-sync.sh >/dev/null 2>&1
